# superviserd.conf
[supervisord]
nodaemon=true

[program:fcgiwrap]
command=/usr/bin/fcgiwrap
autostart=true
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
